# target/linux/zynq/base-files/etc/init.d/post_upgrade
#!/bin/sh /etc/rc.common

START=96
start() {
	. /lib/zynq.sh
	case $(zynq_board_name) in
		dm1-g9|\
		dm1-g19)
			# after successful upgrade delete environment variable to suppress
			# U-Boot auto-recovery process
			fw_setenv upgrade_stage

			# after successful factory reset delete environment variable to
			# allow booting from SD card; it prevents recursive factory when
			# SD card is still inserted in the slot
			if [ x$(fw_printenv -n factory_reset) != x"yes" ] ; then
				fw_setenv first_boot
				fw_setenv sd_boot
				fw_setenv sd_images
			fi
			;;
	esac
}
