# target/linux/zynq/base-files/etc/init.d/post_upgrade
#!/bin/sh /etc/rc.common

START=96
start() {
	. /lib/zynq.sh
	case $(zynq_board_name) in
		dm1-g9|\
		dm1-g19)
			# after successful upgrade delete environment variable to suppress
			# U-Boot auto-recovery process
			fw_setenv upgrade_stage
			;;
	esac
}
